# local.d/000.stop
#!/bin/bash

/usr/bin/rsync -au --delete /var/cache/ccache/ /var/data/ccache/
