# local.d/001.start
#!/bin/bash

/bin/mount /var/cache/ccache/
/usr/bin/rsync -a /var/data/ccache/ /var/cache/ccache/&
